# configs/autostart/blueman-applet.desktop
[Desktop Entry]
Type=Application
Name=Blueman Applet
Exec=blueman-applet
X-GNOME-Autostart-enabled=true
