# configs/autostart/ib-wallpaper-rotator.desktop
[Desktop Entry]
Type=Application
Name=IB Wallpaper Rotator
Exec=/home/ibrahim/.local/bin/ib-wallpaper-rotator
X-GNOME-Autostart-enabled=true
